# clean-pkgs.timer
[Unit]
Description=Flush Old Builds

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
AccuracySec=15m
Persistent=true

[Install]
WantedBy=multi-user.target
